# rootfs/etc/s6-overlay/s6-rc.d/password-display/up
#!/usr/bin/execlineb -P
exec -c /usr/local/bin/print-setup-password.sh


